# quick-reference/bin/crontab
# run five minutes after 8, daily
5 8 * * *       /home/groups/q/qr/qref/htdocs/qref/bin/qrefnew >> $HOME/qrefnew.log 2>&1
	      